# deployment/skynet_backend.service
[Unit]
Description=Gunicorn instance to serve SkyNet Backend
After=network.target

[Service]
# User must match who owns the files (root is fine for VPS, typically www-data is safer but root works here)
User=root
Group=www-data

# Pointing to your specific backend folder on the VPS
WorkingDirectory=/var/www/marketinsights/backend

# Ensure we use the virtual environment created on the server
Environment="PATH=/var/www/marketinsights/backend/venv/bin"

# Runs the app on localhost:8000
ExecStart=/var/www/marketinsights/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000

[Install]
WantedBy=multi-user.target